# edgeoffinality/crontab
30 23 * * * python3 /scripts/main.py && echo "Done executing scripts!"
